# raw/usr/lib/systemd/system/zimaos-status.service
[Unit]
After=casaos-gateway.service
After=casaos-message-bus.service
After=casaos-user-service.service
Description=ZimaOS Chat Service

[Service]
ExecStart=/usr/bin/zimaos-status
Restart=always

[Install]
WantedBy=multi-user.target
